# templates/jujushell.service
[Unit]
Description=Juju shell terminal server

[Service]
ExecStart={{jujushell}} {{jujushell_config}}
User=ubuntu
